# src/config/program-rsyslog.conf
[program:rsyslog]
command = /etc/rc.d/init.d/rsyslog start
