# src/config/program-crond.conf
[program:crond]
command = /etc/rc.d/init.d/crond start
